# publish/blog-publish.path
[Unit]
Description=Watch the Obsidian blog/Publish folder

[Path]
PathChanged=/home/chris/vault/blog/Publish
Unit=blog-publish.service

[Install]
WantedBy=multi-user.target
